# install/local-dash-bind.init
#!/bin/sh /etc/rc.common
#
# Description: Adds the local.brck.com zone to BIND
#

START=99

start() {
  rndc \
    addzone local.brck.com \
    '{type master; file "/etc/bind/db.local.brck.com";};' \
    2>/dev/null
  /etc/init.d/named restart
}
